# Seminar_Around_PHP/DZ_6/docker/supervisor/conf.d/bot.ini
[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/mysite.local/runner -c handle_events_daemon
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/worker