# Seminar_Around_PHP/DZ_6/docker/supervisor/conf.d/sender.ini
[program:sender]
process_name=%(program_name)s_%(process_num)02d
command=php /app/mysite.local/runner -c queue_manager
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/sender